# gateway/server/serialization-server-minimal/src/main/rpm/scriptlet/postinstall
/sbin/chkconfig --add vnet-gateway-serialization-server
